# templates/etc/update-motd.d/99-footer
#!/bin/bash

echo ""
echo -e "${NC} UTC: ${TEXT_COLOR}$(date -u)${NC}, local: ${TEXT_COLOR}$(date)"
echo ""